# smf/manifest/mako-gc-feeder-us.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
    This Source Code Form is subject to the terms of the Mozilla Public
    License, v. 2.0. If a copy of the MPL was not distributed with this
    file, You can obtain one at http://mozilla.org/MPL/2.0/.
-->

<!--
    Copyright (c) 2019, Joyent, Inc.
-->

<service_bundle type="manifest" name="mako-gc-feeder">
    <service name="manta/application/mako_gc_feeder" type="service" version="1">

	<dependency name="network"
		    grouping="require_all"
		    restart_on="error"
		    type="service">
	    <service_fmri value="svc:/network/physical" />
	</dependency>

	<dependency name="filesystem"
		    grouping="require_all"
		    restart_on="error"
		    type="service">
	    <service_fmri value="svc:/system/filesystem/local" />
	</dependency>

	<dependency name="mdata"
		    grouping="require_all"
		    restart_on="none"
		    type="service">
	    <service_fmri value="svc:/smartdc/mdata:execute" />
	</dependency>

	<exec_method type="method"
		     name="start"
			 exec="node --abort-on-uncaught-exception main.js -f etc/config-${index}.json &amp;"
		     timeout_seconds="30">
		<method_context working_directory="/opt/smartdc/mako-gc-feeder" />
	</exec_method>

	<exec_method type="method"
		     name="stop"
		     exec=":kill"
		     timeout_seconds="30" />

        <instance name="mako-gc-feeder-1" enabled="true">
	    <property_group name="mako_gc_feeder" type="application">
		<propval name="index" type="count" value="@@MAKO_GC_FEEDER_INSTANCE_NUM@@" />
            </property_group>
	</instance>

        <instance name="mako-gc-feeder-2" enabled="true">
	    <property_group name="mako_gc_feeder" type="application">
		<propval name="index" type="count" value="@@MAKO_GC_FEEDER_INSTANCE_NUM@@" />
            </property_group>
	</instance>

        <instance name="mako-gc-feeder-3" enabled="true">
	    <property_group name="mako_gc_feeder" type="application">
		<propval name="index" type="count" value="@@MAKO_GC_FEEDER_INSTANCE_NUM@@" />
            </property_group>
	</instance>

        <instance name="mako-gc-feeder-4" enabled="true">
	    <property_group name="mako_gc_feeder" type="application">
		<propval name="index" type="count" value="@@MAKO_GC_FEEDER_INSTANCE_NUM@@" />
            </property_group>
	</instance>

        <instance name="mako-gc-feeder-5" enabled="true">
	    <property_group name="mako_gc_feeder" type="application">
		<propval name="index" type="count" value="@@MAKO_GC_FEEDER_INSTANCE_NUM@@" />
            </property_group>
	</instance>

        <instance name="mako-gc-feeder-6" enabled="true">
	    <property_group name="mako_gc_feeder" type="application">
		<propval name="index" type="count" value="@@MAKO_GC_FEEDER_INSTANCE_NUM@@" />
            </property_group>
	</instance>

        <instance name="mako-gc-feeder-7" enabled="true">
	    <property_group name="mako_gc_feeder" type="application">
		<propval name="index" type="count" value="@@MAKO_GC_FEEDER_INSTANCE_NUM@@" />
            </property_group>
	</instance>

        <instance name="mako-gc-feeder-8" enabled="true">
	    <property_group name="mako_gc_feeder" type="application">
		<propval name="index" type="count" value="@@MAKO_GC_FEEDER_INSTANCE_NUM@@" />
            </property_group>
	</instance>

        <instance name="mako-gc-feeder-9" enabled="true">
	    <property_group name="mako_gc_feeder" type="application">
		<propval name="index" type="count" value="@@MAKO_GC_FEEDER_INSTANCE_NUM@@" />
            </property_group>
	</instance>

        <instance name="mako-gc-feeder-10" enabled="true">
	    <property_group name="mako_gc_feeder" type="application">
		<propval name="index" type="count" value="@@MAKO_GC_FEEDER_INSTANCE_NUM@@" />
            </property_group>
	</instance>

        <instance name="mako-gc-feeder-11" enabled="true">
	    <property_group name="mako_gc_feeder" type="application">
		<propval name="index" type="count" value="@@MAKO_GC_FEEDER_INSTANCE_NUM@@" />
            </property_group>
	</instance>

        <instance name="mako-gc-feeder-12" enabled="true">
	    <property_group name="mako_gc_feeder" type="application">
		<propval name="index" type="count" value="@@MAKO_GC_FEEDER_INSTANCE_NUM@@" />
            </property_group>
	</instance>

	<stability value='Unstable' />

	<template>
	    <common_name>
		<loctext xml:lang="C">Joyent Manta Storage Service Public API</loctext>
	    </common_name>
	</template>

    </service>
</service_bundle>
