# smf/manifest/mako-gc-feeder-eu.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
    This Source Code Form is subject to the terms of the Mozilla Public
    License, v. 2.0. If a copy of the MPL was not distributed with this
    file, You can obtain one at http://mozilla.org/MPL/2.0/.
-->

<!--
    Copyright (c) 2019, Joyent, Inc.
-->

<service_bundle type="manifest" name="mako-gc-feeder">
    <service name="manta/application/mako_gc_feeder" type="service" version="1">

	<dependency name="network"
		    grouping="require_all"
		    restart_on="error"
		    type="service">
	    <service_fmri value="svc:/network/physical" />
	</dependency>

	<dependency name="filesystem"
		    grouping="require_all"
		    restart_on="error"
		    type="service">
	    <service_fmri value="svc:/system/filesystem/local" />
	</dependency>

	<dependency name="mdata"
		    grouping="require_all"
		    restart_on="none"
		    type="service">
	    <service_fmri value="svc:/smartdc/mdata:execute" />
	</dependency>

	<exec_method type="method"
		     name="start"
		     exec="/root/mako-gc-feeder-current/tools/node/bin/node --abort-on-uncaught-exception main.js -f etc/config-%i.json &amp;"
		     timeout_seconds="30">
		<method_context working_directory="/root/mako-gc-feeder-current" />
	</exec_method>

	<exec_method type="method"
		     name="stop"
		     exec=":kill"
		     timeout_seconds="30" />

        <instance name="mako-gc-feeder-0"  enabled="false" />
        <instance name="mako-gc-feeder-1"  enabled="false" />
        <instance name="mako-gc-feeder-2"  enabled="false" />
        <instance name="mako-gc-feeder-3"  enabled="false" />
        <instance name="mako-gc-feeder-4"  enabled="false" />
        <instance name="mako-gc-feeder-5"  enabled="false" />
        <instance name="mako-gc-feeder-6"  enabled="false" />
        <instance name="mako-gc-feeder-7"  enabled="false" />
        <instance name="mako-gc-feeder-8"  enabled="false" />
        <instance name="mako-gc-feeder-9"  enabled="false" />
        <instance name="mako-gc-feeder-10" enabled="false" />
        <instance name="mako-gc-feeder-11" enabled="false" />
        <instance name="mako-gc-feeder-12" enabled="false" />
        <instance name="mako-gc-feeder-13" enabled="false" />
        <instance name="mako-gc-feeder-14" enabled="false" />

	<stability value='Unstable' />

	<template>
	    <common_name>
		<loctext xml:lang="C">Joyent Manta Storage Service Public API</loctext>
	    </common_name>
	</template>

    </service>
</service_bundle>
